# deploy/aiproxy.service
[Unit]
Description=aiboard CORS proxy -> https://api.zenoid.space
After=network.target

[Service]
Type=simple
User=aiproxy
Group=aiproxy
WorkingDirectory=/srv/aiboard
ExecStart=/usr/bin/python3 /srv/aiboard/proxy.py 8787 https://api.zenoid.space
Restart=always
RestartSec=2

# Sandboxing — прокси читает только свой каталог
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ReadOnlyPaths=/srv/aiboard

[Install]
WantedBy=multi-user.target
